# pi-config/etc/rc.local
#!/bin/sh -e
# Pi 上 /etc/rc.local 的内容快照
# 现已为空（sugar-wifi-conf 由 sugar-wifi-config.service 接管）

_IP=$(hostname -I) || true
if [ "$_IP" ]; then
  printf "My IP address is %s\n" "$_IP"
fi

exit 0
